# cli/vaultkeeperd.service
[Unit]
Description=VaultKeeper P2P Storage Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vaultkeeper
Group=vaultkeeper
ExecStart=/usr/local/bin/vaultkeeperd daemon --data-dir /var/lib/vaultkeeper
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vaultkeeper
PrivateTmp=true

# Sandboxing — restrict syscalls to system-service set
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
CapabilityBoundingSet=

# Resource limits
MemoryMax=2G
CPUQuota=50%

# Additional hardening
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=multi-user.target
